# assets/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info

[program:smbd]
command=/usr/sbin/smbd --daemon --foreground --log-stdout
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nmbd]
command=/usr/sbin/nmbd --daemon --foreground --log-stdout
redirect_stderr=true

[program:winbind]
command=/usr/sbin/winbindd --daemon --foreground --stdout
redirect_stderr=true